# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true

[program:coolwsd]
command=/usr/bin/coolwsd --no-daemon
autorestart=true
